# mappings/Azure/AzureSentinelHunting-51-100.yaml
version: 1
ATT&CK version: 8.1
creation date: 04/23/2021
name: Azure Sentinel Hunting 51-100
contact: ctid@mitre-engenuity.org
organization: Center for Threat Informed Defense (CTID)
platform: Azure
description: Out of the box Azure Sentinel Hunting (from the queries)
techniques:
  - id: T1105
    name: Ingress Tool Transfer
    technique-scores:
      - category: Detect
        value: Minimal
        comments: >-
          This control may detect "Crypto currency miners EXECVE" which may indicate download of
          crypto currency miners through EXECVE.
  - id: T1059
    name: Command and Scripting Interpreter
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1059.003
            name: Windows Command Shell
          - id: T1059.007
            name: JavaScript/JScript
          - id: T1059.005
            name: Visual Basic
          - id: T1059.006
            name: Python
        scores:
          - category: Detect
            value: Minimal
            comments: >-
              This control may detect malicious scripting by providing a "Cscript script daily
              summary breakdown". This control may also detect "Hosts running a rare process with
              commandline".
  - id: T1071
    name: Application Layer Protocol
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1071.004
            name: DNS
        scores:
          - category: Detect
            value: Partial
            comments: >-
              This control may detect "DNS - domain anomalous lookup increase", "DNS Full Name
              anomalous lookup increase", or "DNS lookups for commonly abused TLDs", which shows an
              increase in domain lookups for a client IP and indicate malicious traffic or
              exfiltration of sensitive data.
  - id: T1053
    name: Scheduled Task/Job
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1053.003
            name: Cron
        scores:
          - category: Detect
            value: Partial
            comments: This control may detect "Editing Linux scheduled tasks through Crontab".
  - id: T1087
    name: Account Discovery
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1087.002
            name: Domain Account
          - id: T1087.001
            name: Local Account
        scores:
          - category: Detect
            value: Minimal
            comments: >-
              This control may detect "Enumeration of users and groups" through use of the 'net'
              tool.
  - id: T1069
    name: Permission Groups Discovery
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1069.002
            name: Domain Groups
          - id: T1069.001
            name: Local Groups
        scores:
          - category: Detect
            value: Minimal
            comments: >-
              This control may detect "Enumeration of users and groups" through use of the 'net'
              tool.
  - id: T1036
    name: Masquerading
    technique-scores:
      - category: Detect
        value: Minimal
        comments: >-
          This control may detect "Exes with double file extension and access summary" which may
          indicate malicious executable files that have been hidden as other file types.
  - id: T1136
    name: Create Account
    technique-scores:
      - category: Detect
        value: Minimal
        comments: >-
          This control may detect "External user added and removed in short timeframe" or "External
          user from a new organisation added" which may indicate establishment of persistence. These
          queries only work with Teams user accounts. This control may also detect "Group added to
          Built in Domain Local or Global Group" which may indicate privilege escalation.
  - id: T1110
    name: Brute Force
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1110.004
            name: Credential Stuffing
          - id: T1110.003
            name: Password Spraying
          - id: T1110.001
            name: Password Guessing
        scores:
          - category: Detect
            value: Partial
            comments: >-
              This control may detect "Failed attempt to access Azure Portal", "Failed Login Attempt
              by Expired account", "Failed Logon Attempts on SQL Server", Failed Logon on SQL Server
              from Same IPAddress in Short time Span", "Failed service logon attempt by user account
              with available AuditData", "Login attempt by Blocked MFA user", or "Login spike with
              increase failure rate".
  - id: T1018
    name: Remote System Discovery
    technique-scores:
      - category: Detect
        value: Minimal
        comments: >-
          This control may detect "High reverse DNS count by host" which may indicate remote system
          discovery.
  - id: T1114
    name: Email Collection
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1114.002
            name: Remote Email Collection
        scores:
          - category: Detect
            value: Minimal
            comments: >-
              This control may detect "Host Exporting Mailbox and Removing Export" which may
              indicate exfiltration of data from Exchange servers.
      - sub-techniques:
          - id: T1114.003
            name: Email Forwarding Rule
        scores:
          - category: Detect
            value: Minimal
            comments: >-
              This control may detect "Mail redirect via ExO transport rule". This is specific to
              Exchange servers only.
  - id: T1078
    name: Valid Accounts
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1078.003
            name: Local Accounts
          - id: T1078.002
            name: Domain Accounts
        scores:
          - category: Detect
            value: Minimal
            comments: >-
              This control may detect "Hosts with new logons" which may indicate lateral movement
              using existing valid accounts. This control may also detect "Inactive or new account
              signins" or "Long lookback User Account Created and Deleted within 10mins".
references:
  - 'https://docs.microsoft.com/en-us/azure/sentinel/hunting'
