# mappings/Azure/FileIntegrityMonitoring.yaml
version: 1
ATT&CK version: 8.1
creation date: 1/21/2021
name: File Integrity Monitoring
contact: ctid@mitre-engenuity.org
organization: Center for Threat Informed Defense (CTID)
platform: Azure
tags:
  - Azure Security Center
  - Azure Security Center Recommendation
  - Azure Defender
  - Azure Defender for Servers
  - File system
  - Registry
  - Windows
  - Linux
description: >
  File integrity monitoring (FIM), also known as change monitoring, examines  operating system
  files, Windows registries, application software, Linux  system files, and more, for changes that
  might indicate an attack. File Integrity Monitoring (FIM) informs you when changes occur to
  sensitive  areas in your resources, so you can investigate and address unauthorized  activity. 
techniques:
  - id: T1053
    name: Scheduled Task/Job
    technique-scores:
      - category: Detect
        value: Significant
    sub-techniques-scores:
      - sub-techniques:
          - id: T1053.001
            name: At (Linux)
          - id: T1053.002
            name: At (Windows)
          - id: T1053.003
            name: Cron
          - id: T1053.005
            name: Scheduled Task
          - id: T1053.006
            name: Systemd Timers
        scores:
          - category: Detect
            value: Significant
  - id: T1098
    name: Account Manipulation
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1098.004
            name: SSH Authorized Keys
        scores:
          - category: Detect
            value: Significant
  - id: T1547
    name: Boot or Logon Autostart Execution
    technique-scores:
      - category: Detect
        value: Significant
    sub-techniques-scores:
      - sub-techniques:
          - id: T1547.001
            name: Registry Run Keys / Startup Folder
          - id: T1547.002
            name: Authentication Package
          - id: T1547.003
            name: Time Providers
          - id: T1547.004
            name: Winlogon Helper DLL
          - id: T1547.005
            name: Security Support Provider
          - id: T1547.006
            name: Kernel Modules and Extensions
          - id: T1547.008
            name: LSASS Driver
          - id: T1547.009
            name: Shortcut Modification
          - id: T1547.010
            name: Port Monitors
          - id: T1547.012
            name: Print Processors
        scores:
          - category: Detect
            value: Significant
  - id: T1037
    name: Boot or Logon Initialization Scripts
    technique-scores:
      - category: Detect
        value: Significant
    sub-techniques-scores:
      - sub-techniques:
          - id: T1037.001
            name: Logon Script (Windows)
          - id: T1037.003
            name: Network Logon Script
        scores:
          - category: Detect
            value: Significant
  - id: T1543
    name: Create or Modify System Process
    technique-scores:
      - category: Detect
        value: Significant
    sub-techniques-scores:
      - sub-techniques:
          - id: T1543.002
            name: Systemd Service
          - id: T1543.003
            name: Windows Service
        scores:
          - category: Detect
            value: Significant
  - id: T1546
    name: Event Triggered Execution
    technique-scores:
      - category: Detect
        value: Partial
        comments: >
          The detection score for this technique was assessed as Partial because it doesn't detect
          some of the sub-techniques of this technique such as Windows Management Instrumentation
          Event Subscription (WMI) and Trap sub-techniques.  Additionally for  some sub-techniques,
          this control can be noisy.
    sub-techniques-scores:
      - sub-techniques:
          - id: T1546.001
            name: Change Default File Association
          - id: T1546.002
            name: Screensaver
          - id: T1546.004
            name: .bash_profile and .bashrc
          - id: T1546.007
            name: Netsh Helper DLL
          - id: T1546.008
            name: Accessibility Features
          - id: T1546.009
            name: AppCert DLLs
          - id: T1546.011
            name: Application Shimming
          - id: T1546.012
            name: Image File Execution Options Injection
          - id: T1546.013
            name: PowerShell Profile
        scores:
          - category: Detect
            value: Significant
      - sub-techniques:
          - id: T1546.010
            name: AppInit DLLs
          - id: T1546.015
            name: Component Object Model Hijacking
        scores:
          - category: Detect
            value: Minimal
            comments: >
              The detection score for this group of sub-techniques is assessed as Minimal due to the
              accuracy component of the score.  The registry keys which are modified as a result of
              these sub-techniques can change frequently or are  too numerous to monitor and
              therefore can result in significant amount of false positives.
  - id: T1574
    name: Hijack Execution Flow
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1574.006
            name: LD_PRELOAD
        scores:
          - category: Detect
            value: Partial
  - id: T1137
    name: Office Application Startup
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1137.002
            name: Office Test
        scores:
          - category: Detect
            value: Significant
  - id: T1548
    name: Abuse Elevation Control Mechanism
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1548.002
            name: Bypass User Account Control
        scores:
          - category: Detect
            value: Minimal
            comments: >-
              Some UAC bypass methods rely on modifying specific, user-accessible Registry settings
              that can be monitored using this control.  Overall, there are numerous other bypass
              methods that do not result in Registry modification that this control will not be
              effect in detection resulting in a low detection coverage factor.
      - sub-techniques:
          - id: T1548.003
            name: Sudo and Sudo Caching
        scores:
          - category: Detect
            value: Significant
  - id: T1556
    name: Modify Authentication Process
    technique-scores:
      - category: Detect
        value: Partial
        comments: >-
          This control is effective for detecting the Registry and file system artifacts that are
          generated during the execution of some variations of this technique while minimizing false
          positives due to the locations being monitored changing infrequently (e.g. /etc/pam.d/). 
    sub-techniques-scores:
      - sub-techniques:
          - id: T1556.002
            name: Password Filter DLL
        scores:
          - category: Detect
            value: Significant
            comments: >-
              The Registry key used to register a Password Filter DLL can be monitored for changes
              using this control providing substantial coverage of this sub-technique.  This key
              should not change often and therefore false positives should be minimal.
      - sub-techniques:
          - id: T1556.003
            name: Pluggable Authentication Modules
        scores:
          - category: Detect
            value: Significant
            comments: >-
              The PAM configuration and module paths (/etc/pam.d/) can be monitored for changes
              using this control.  The files in this path should not change often and therefore
              false positives should be minimal.
  - id: T1003
    name: OS Credential Dumping
    technique-scores:
      - category: Detect
        value: Minimal
        comments: >-
          Most credential dumping operations do not require modifying resources that can be detected
          by this control (i.e. Registry and File system) and therefore its coverage is minimal.
    sub-techniques-scores:
      - sub-techniques:
          - id: T1003.001
            name: LSASS Memory
        scores:
          - category: Detect
            value: Partial
            comments: >-
              This control can be used to detect the Windows Security Support Provider (SSP) DLLs
              variation of this sub-technique by monitoring the Registry keys used to register these
              DLLs.  These keys should change infrequently and therefore false positives should be
              minimal. 
  - id: T1222
    name: File and Directory Permissions Modification
    technique-scores:
      - category: Detect
        value: Partial
        comments: >-
          This control can detect file and directory permissions modifications for Windows but there
          is no indication it can do so for Linux platforms.
    sub-techniques-scores:
      - sub-techniques:
          - id: T1222.001
            name: Windows File and Directory Permissions Modification
        scores:
          - category: Detect
            value: Significant
            comments: >-
              This control can detect changes to the Discretionary Access Control Lists (DACLs) of
              Windows files and can be used to detect modifications to sensitive directories and
              files that shouldn't change frequently.
  - id: T1562
    name: Impair Defenses
    technique-scores:
      - category: Detect
        value: Minimal
        comments: 'Due to low detection coverage, this technique is scored as minimal.'
    sub-techniques-scores:
      - sub-techniques:
          - id: T1562.001
            name: Disable or Modify Tools
        scores:
          - category: Detect
            value: Minimal
            comments: >-
              This control can be used to monitor Registry keys related to security software or
              event logging processes that can detect when an adversary attempts to disable these
              tools via modifying or deleting Registry keys.  A majority of the cited procedure
              examples for this sub-technique are related to killing security processes rather than
              modifying the Registry, and therefore the detection coverage for this control is low.
      - sub-techniques:
          - id: T1562.004
            name: Disable or Modify System Firewall
          - id: T1562.006
            name: Indicator Blocking
        scores:
          - category: Detect
            value: Minimal
            comments: >-
              There are numerous ways depending on the operating system that these sub-techniques
              can be accomplished.  Monitoring the Windows Registry is one way depending on the
              procedure chosen to implement the sub-technique and therefore the overall coverage is
              low.
  - id: T1553
    name: Subvert Trust Controls
    technique-scores:
      - category: Detect
        value: Partial
        comments: >-
          This control can be used to detect a subset of this technique's sub-techniques while
          minimizing the false positive rate.
    sub-techniques-scores:
      - sub-techniques:
          - id: T1553.003
            name: SIP and Trust Provider Hijacking
        scores:
          - category: Detect
            value: Partial
            comments: >-
              This control can detect modifications made to the Registry keys used to register
              Windows Subject Interface Packages (SIPs) .  Because this sub-technique can be
              accomplished without modifying the Registry via DLL Search Order Hijacking, it has
              been scored as Partial.  The related Registry keys should change often and therefore
              the false positive rate should be minimal.
      - sub-techniques:
          - id: T1553.004
            name: Install Root Certificate
        scores:
          - category: Detect
            value: Significant
            comments: >-
              This control can be used to detect when the system root certificates has changed by
              detecting the corresponding Registry or File system modifications that occur as a
              result.  These root certificates should not change often and therefore the false
              positive rate is minimal.
comments: >
  These techniques included in this mapping result in Windows Registry or file system artifacts
  being created or modified which can be detected by this control.  

  The detection score for most techniques included in this mapping was scored as Significant  and
  where there are exceptions, comments have been provided. This Significant score assessment  was
  due to the following factors: Coverage - (High) The control was able to detect most of the
  sub-techniques, references and procedure examples of the mapped techniques. Accuracy - (High)
  Although this control does not include built-in intelligence to minimize  the false positive rate,
  the specific artifacts generated by the techniques in this mapping  do not change frequently and
  therefore the potential for a high false-positive is reduced.  Temporal - (Medium) This control at
  worst scans for changes on an hourly basis.
references:
  - 'https://docs.microsoft.com/en-us/azure/security-center/security-center-file-integrity-monitoring'
